# packaging/systemd/waired-agent.service
[Unit]
Description=Waired agent
Documentation=https://github.com/waired-ai/waired
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/waired/agent.env
ExecStart=/usr/bin/waired-agent --state-dir /var/lib/waired
User=waired
Group=waired
# Restart=always (not on-failure): a management-API-driven reconfigure exits
# cleanly via SIGTERM (exit 0), and the proxy enable/disable + preferred-model
# flows rely on systemd bringing the agent back after that clean exit.
# on-failure would treat the SIGTERM exit as success and NOT restart. Matches
# the CLI-rendered unit (internal/platform/service/service_linux.go).
Restart=always
RestartSec=5s
# Exit code 17 is the agent's "restart me" signal (preferred-model
# switch via the management API, issue #347): counted as success so
# the unit never enters the failed state, but force-restarted so the
# new selection takes effect. Crashes restart via on-failure; plain
# exit 0 and `systemctl stop` stay down.
SuccessExitStatus=17
RestartForceExitStatus=17

ProtectSystem=strict
ReadWritePaths=/var/lib/waired
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
